# pkgs/logrotate/logrotate.cron
MAILTO=""
30 2 * * * root %BASE%/sbin/logrotate %BASE%/etc/logrotate.conf > %BASE%/var/log/logrotate.log 2>&1

